# shadowsocks.conf
[program:shadowsocks]
command=ss-server -c /etc/shadowsocks.json
autorestart=true
user=nobody
